# 麒麟系统安装包/用户工具安装包/source/x86/deb/lm-gm3000user.i386/etc/xdg/autostart/lm-gm3000user.desktop
[Desktop Entry]
Name=GM3000User
GenericName=GM3000User
Comment=GM3000User
Comment[zh_CN]=GM3000User
Exec=/opt/lm-gm3000user/script/AppRunAuto.sh
Terminal=false
Type=Application
Icon=/opt/lm-gm3000user/logo/logo.png
Categories=Utility;Application;
StartupNotify=false

