# UOS系统安装包/超级插件安装包/source/飞腾FT1500/com.longmai.mpluginv4/opt/apps/com.longmai.mpluginv4/entries/autostart/com.longmai.mpluginv4.desktop
[Desktop Entry]
Exec=/opt/apps/com.longmai.mpluginv4/files/script/run.sh
Type=Application
Terminal=false
Name=mPluginV4
GenericName=mPluginV4